# crontab
* * * * * /home/German/Desktop/ger/code/projects/scripts/switchwallpapers.sh
